# server/mediamtx.service
[Unit]
Description=MediaMTX RTSP Server
After=network.target

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/mediamtx /etc/mediamtx.yml
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
